# rc/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

ipset create user bitmap:ip,mac range 202.141.166.0/24 timeout 604800

#ipset add user 202.141.166.2,38:bc:1a:00:ac:e6

/etc/rc.d/rc.firewall
