# rc/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local


SUBNET=192.168.0.0/20

#logined users
ipset create user bitmap:ip,mac range $SUBNET timeout 604800

#new users (could not do auto login)
ipset create newuser bitmap:ip,mac range $SUBNET timeout 604800

#ipset add user 202.141.166.2,38:bc:1a:00:ac:e6

/etc/rc.d/rc.firewall

echo 120 > /proc/sys/net/ipv4/neigh/default/gc_stale_time
echo 2048 > /proc/sys/net/ipv4/neigh/default/gc_thresh1
echo 4096 > /proc/sys/net/ipv4/neigh/default/gc_thresh2
echo 8192 > /proc/sys/net/ipv4/neigh/default/gc_thresh3

echo 100000 > /proc/sys/net/netfilter/nf_conntrack_max

echo 600 > /proc/sys/net/netfilter/nf_conntrack_tcp_timeout_established
echo 10 > /proc/sys/net/netfilter/nf_conntrack_udp_timeout
echo 120 > /proc/sys/net/netfilter/nf_conntrack_udp_timeout_stream


